# crontab
*/1 * * * * /bin/sh /gohugo >/dev/stdout 2>&1